# systemd/exit-gate-minimal.service
[Unit]
Description=Exit Gate - Linux Application Firewall Daemon
Documentation=https://github.com/exit-gate/exit-gate
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/exit-gate-daemon --config /etc/exit-gate/config.toml --foreground
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=exit-gate

# Security settings - temporarily relaxed for testing
NoNewPrivileges=false
#PrivateTmp=true
#ProtectSystem=strict
#ProtectHome=true
#ReadWritePaths=/var/lib/exit-gate /var/run/exit-gate /run/exit-gate

# Capabilities needed for eBPF
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_NET_ADMIN CAP_SYS_ADMIN CAP_DAC_OVERRIDE

# Resource limits
LimitNOFILE=65536
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
